# docker/supervisor-app.conf
[program:uwsgi]
command = /usr/sbin/uwsgi --ini /Savior/docker/savior.ini
user = nobody
stopsignal = QUIT

[program:nginx]
command = /usr/sbin/nginx -g 'daemon off;'
